# pltn_video_display/pltn_video_display.service
[Unit]
Description=PLTN Video Display System
Documentation=https://github.com/Sanssin/pkm-simulator-PLTN
After=network.target graphical.target pkm-simulator.service
Wants=graphical.target
Requires=pkm-simulator.service

[Service]
Type=simple
User=pkm
Group=pkm
WorkingDirectory=/home/pkm/pkm-simulator-PLTN/pltn_video_display
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pkm/.Xauthority"

# Start video display in fullscreen mode
ExecStart=/usr/bin/python3 /home/pkm/pkm-simulator-PLTN/pltn_video_display/video_display_app.py

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pltn_video_display

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical.target
