# pltn_video_display/pltn_video_display_nouser.service
[Unit]
Description=PLTN Video Display System
Documentation=https://github.com/Sanssin/pkm-simulator-PLTN
After=network.target graphical.target pkm-simulator.service
Wants=graphical.target
Requires=pkm-simulator.service

[Service]
Type=simple
# Run as root to avoid user permission issues
WorkingDirectory=/home/pi/pkm-simulator-PLTN/pltn_video_display
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
Environment="HOME=/home/pi"

# Start video display in fullscreen mode
ExecStart=/usr/bin/python /home/pi/pkm-simulator-PLTN/pltn_video_display/video_display_app.py

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pltn_video_display

[Install]
WantedBy=graphical.target
